# config/etc/services.d/proxysql/run
#!/bin/sh
proxysql -f -c /etc/proxysql/proxysql.cnf --reload
